# rootfs/etc/systemd/system/fail2ban.service.d/override.conf
[Service]
ExecStartPre=/warewulf/sbin/execstartpre.%N.sh
Restart=no

